# lib/engineyard-hudson/cli/server/templates/cookbooks/hudson_master/templates/default/init.sh.erb
#!/sbin/runscript

# Start/stop script for hudson

PIDFILE=<%= @pid %>
HUDSON_HOME=<%= @home %>
COMMAND=/usr/bin/java
ARGUMENTS="-Djava.io.tmpdir=${HUDSON_HOME}/tmp -jar ${HUDSON_HOME}/hudson.war --webroot=${HUDSON_HOME}/war --httpPort=<%= @port %> ajp13ListenAddress=127.0.0.1"
RUN_AS=<%= @user %>
NAME=Hudson-CI

start() {
  ebegin "Starting $NAME"
  export HUDSON_HOME=${HUDSON_HOME}
  start-stop-daemon -b --start --make-pidfile --pidfile $PIDFILE --chuid $RUN_AS --exec $COMMAND -- $ARGUMENTS &>$HUDSON_HOME/logs/hudson_log
  eend $?
}

stop() {
  ebegin "Stopping $NAME"
  start-stop-daemon --stop --quiet --pidfile $PIDFILE
  rm -f $PIDFILE
  eend $?
}